# supervisord.conf
[supervisord]
nodaemon = true
user = root

[program:telldusmq]
command = /go/bin/telldusmq

[program:telldusd]
command = telldusd --nodaemon